# workspace/crontab/laradock
* * * * * laradock php /var/www/breedregistry/artisan schedule:run >> /dev/null 2>&1
